# html5-desktop/app/conf.d/npm.conf
[program:npm]
command=npm --prefix /home/api start
autorestart=true